# ops/launchd/com.mailvec.indexer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mailvec.indexer</string>
    <key>ProgramArguments</key>
    <array>
        <string>__DOTNET__</string>
        <string>__INSTALL_PREFIX__/Mailvec.Indexer.dll</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__INSTALL_PREFIX__</string>
    <!-- Settings that vary per install (DB path, Maildir root, Ollama URL,
         Fastmail account id) live in the shared file at
         ~/Library/Application Support/Mailvec/appsettings.Local.json — the
         single source of truth read by every Mailvec binary. ops/install.sh
         writes it. Only env vars genuinely launchd-specific stay here. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>MAILVEC_LAUNCHD</key>
        <string>1</string>
        <key>MAILVEC_LOG_DIR</key>
        <string>__LOG_DIR__</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <!-- Serilog writes the rolling primary log to MAILVEC_LOG_DIR. These two
         paths only capture pre-Serilog startup output and unhandled native
         stderr (panics, runtime asserts). They grow slowly; no rotation needed. -->
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/mailvec-indexer.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/mailvec-indexer.launchd.log</string>
</dict>
</plist>
