# ops/launchd/com.mailvec.mbsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mailvec.mbsync</string>
    <key>ProgramArguments</key>
    <array>
        <string>__MBSYNC__</string>
        <string>-c</string>
        <string>__MBSYNCRC__</string>
        <string>-a</string>
    </array>
    <!-- 600s (10 min) instead of 300s. mbsync runs occasionally take longer
         than 5 minutes when pulling a backlog, and a too-tight schedule
         causes the next run to hit flock() on .mbsyncstate and exit with
         "channel is locked". 92% of runs were succeeding at 300s, but the
         8% lock failures were producing noisy tray banners. 10 min still
         gives near-real-time mail delivery on a personal mailbox. -->
    <key>StartInterval</key>
    <integer>600</integer>
    <key>RunAtLoad</key>
    <true/>
    <!-- mbsync is the only non-.NET service, so Serilog doesn't apply. Its
         output goes through launchd-captured paths instead. Volume is small
         (few lines per 5-minute sync); rotation isn't worth the complexity. -->
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/mailvec-mbsync.out.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/mailvec-mbsync.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
