# ops/launchd/com.mailvec.mbsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mailvec.mbsync</string>
    <key>ProgramArguments</key>
    <array>
        <string>__MBSYNC__</string>
        <string>-c</string>
        <string>__MBSYNCRC__</string>
        <string>-a</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <!-- mbsync is the only non-.NET service, so Serilog doesn't apply. Its
         output goes through launchd-captured paths instead. Volume is small
         (few lines per 5-minute sync); rotation isn't worth the complexity. -->
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/mailvec-mbsync.out.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/mailvec-mbsync.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
