# ops/launchd/com.mailvec.embedder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mailvec.embedder</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/share/dotnet/dotnet</string>
        <string>__INSTALL_PREFIX__/Mailvec.Embedder.dll</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__INSTALL_PREFIX__</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/mailvec-embedder.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/mailvec-embedder.err.log</string>
</dict>
</plist>
